# tag-mac/Library/LaunchAgents/setenv.zsh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>setenv.zsh</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/launchctl</string>
      <string>setenv</string>
      <string>ZDOTDIR</string>
      <string>$HOME/.config/zsh</string>
    </array>
    <key>RunAtLoad</key>
    <true />
    <key>StandardErrorPath</key>
    <string>$HOME/Library/Logs/setenv.pg.start.log</string>
    <key>StandardOutPath</key>
    <string>$HOME/Library/Logs/setenv.pg.start.log</string>
    <key>WorkingDirectory</key>
    <string>/tmp</string>
  </dict>
</plist>
